# crates/uaa-control/systemd/uaa-control.service
# file: crates/uaa-control/systemd/uaa-control.service
# version: 1.0.0
# guid: 5c164116-d4a3-4715-b230-aa9e1a1e88e3
# last-edited: 2026-07-10
#
# Service unit for the uaa-control daemon (spec C3). Paired with uaa-control.socket so
# the :25000 machine plane is socket-activated (Decision 24). The :7443/:7444/:8443
# planes are bound by the process itself.
#
# DEPLOY IS BUCKET-3 (human) WORK: docs/reference artifact only. An operator installs
# this on 172.16.2.30; no agent deploys it. Paths/user below are cutover placeholders.

[Unit]
Description=uaa constellation control daemon (spec C3)
Requires=uaa-control.socket
After=uaa-control.socket network-online.target

[Service]
Type=exec
ExecStart=/usr/local/bin/uaa-control serve
Restart=on-failure
# StateDirectory gives /var/lib/uaa for the registry snapshot + WAL (Decision 4).
StateDirectory=uaa
User=uaa
Group=uaa

[Install]
WantedBy=multi-user.target
